# build/deb/DEBIAN/postinst
#!/bin/bash

systemctl enable gnmi.service

systemctl start gnmi.service

