# build/deb/lib/systemd/system/gnmi.service
[Unit]
Description=The gNMI Server
After=docker.service

[Service]
Type=Idle
ExecStart=/usr/local/bin/gnmi run --address 0.0.0.0 --port 5002 -v --config /var/run/redis/sonic-db/database_config.json --username admin --password YourPaSsWoRd

Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target

